# baseos/_/alpine/3.12/vminit-svc
#!/sbin/openrc-run
command="/usr/bin/vminit"
command_args=""
pidfile="/run/vminit.pid"